# worker/supervisord.conf
[supervisord]
user=root  # 消除root用户警告
logfile=/tmp/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=true

minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; [program:django]
; command=python manage.py runserver 0.0.0.0:8000
; directory=/app
; user=root
; autostart=true
; autorestart=true
; redirect_stderr=true
; # 日志同时输出到文件和控制台
; stdout_logfile=/var/log/django_stdout.log
; stdout_logfile_maxbytes=50MB
; stdout_logfile_backups=3
; stderr_logfile=/var/log/django_stderr.log
; stderr_logfile_maxbytes=50MB
; stderr_logfile_backups=3

[program:celery_worker]
command=celery -A config.celery worker --loglevel=info --concurrency=4
 directory=/app
user=root
autostart=true
autorestart=true
redirect_stderr=true
# 日志同时输出到文件和控制台
stdout_logfile=/var/log/celery_worker_stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_logfile=/var/log/celery_worker_stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=3
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr


[program:celery_beat]
command=celery -A config.celery beat --loglevel=info
 directory=/app
user=root
autostart=true
autorestart=true
redirect_stderr=true
# 日志同时输出到文件和控制台
stdout_logfile=/var/log/celery_beat_stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
stderr_logfile=/var/log/celery_beat_stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=3
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr